# startup/nebulas-privatenet.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/nebulas-private.sh

[Install]
WantedBy=default.target